# backend/juvantia-backend.service
[Unit]
Description=Juvantia Backend API Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/juvantia/backend
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=juvantia-backend

[Install]
WantedBy=multi-user.target
